# 09_role-collections/collections/ansible_collections/fedora/linux_system_roles/roles/trustee_server/files/quadlet/trustee-as.container
# Trustee Attestation Service (AS)
# Verifies TEE evidence and generates attestation tokens
#
# AS depends on RVPS for reference values and is required by KBS for attestation.
# It acts as the "Verifier" in the RATS architecture.

[Unit]
Description=Trustee Attestation Service
Documentation=https://docs.redhat.com/en/documentation/openshift_sandboxed_containers/1.10/html/deploying_trustee/

# Start after RVPS is available
After=network-online.target trustee-rvps.service
Wants=network-online.target
Requires=trustee-rvps.service

[Container]
# Container image - uses Red Hat registry
# For development/testing, override with: Image=quay.io/confidential-containers/as:latest
Image=quay.io/linux-system-roles/trustee-as:0.15.0
Pod=trustee.pod

# Container name for DNS resolution on the trustee network
ContainerName=trustee-as

# Configuration mount (read-only)
Volume=/etc/trustee/as:/etc/as:ro,Z

# Data persistence
Volume=as-data.volume:/var/lib/as:Z

# Environment variables
Environment=RUST_LOG=info
Environment=AS_CONFIG_FILE=/etc/as/config.json

# Security hardening
NoNewPrivileges=true

# Drop all capabilities
DropCapability=ALL

# Health check
HealthCmd=sh -c 'ss -tlnp | grep -q 50004 || exit 1'
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=15s

# Resource limits (adjust based on workload)
# Memory=512M

[Service]
# Restart policy
Restart=on-failure
RestartSec=10s

# Startup/shutdown timeouts
TimeoutStartSec=90s
TimeoutStopSec=30s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trustee-as

[Install]
WantedBy=multi-user.target default.target
