# 09_role-collections/collections/ansible_collections/fedora/linux_system_roles/roles/trustee_server/templates/secret_registration_server.service.j2
# SPDX-License-Identifier: MIT
{{ ansible_managed | comment }}
{{ "system_role:trustee_server" | comment(prefix="", postfix="") }}

[Unit]
Description=Secret Registration Server - attestation-backed disk key registration
Documentation=https://confidentialcontainers.org/docs/attestation/
After=network-online.target
Wants=network-online.target
# Start after Trustee pod if it exists
After=trustee-pod.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/secret_registration_server.py
Restart=on-failure
RestartSec=5s
# Required for KBS and attestation access
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/trustee/version.env

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
