# 09_role-collections/collections/ansible_collections/fedora/linux_system_roles/roles/trustee_client/templates/secret_registration_client.service.j2
# SPDX-License-Identifier: MIT
{{ ansible_managed | comment }}
{{ "system_role:trustee_client" | comment(prefix="", postfix="") }}

[Unit]
Description=Secret Registration Client - register with server and fetch disk encryption key
Documentation=https://confidentialcontainers.org/docs/attestation/
After=network-online.target
Wants=network-online.target
# Run after Trustee GC pod (CDH)
After=trustee-gc-pod.service
Wants=trustee-gc-pod.service

[Service]
Type=oneshot
RemainAfterExit=yes
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60
ExecStart=/usr/local/bin/secret_registration_client.sh --mount-disk
EnvironmentFile=-/etc/secret-registration/client.env

[Install]
WantedBy=multi-user.target
